# infra/systemd/flashback-backend.service
[Unit]
Description=Flashback Restore - Backend FastAPI
After=network.target
PartOf=flashback-arq-worker.service

[Service]
Type=simple
User=flashback
Group=flashback
WorkingDirectory=/opt/flashback-restore-monorepo/backend
ExecStart=/opt/flashback-restore-monorepo/backend/.venv/bin/python -m uvicorn app.main:app --host 0.0.0.0 --port 8000 --loop asyncio
Restart=always
RestartSec=5
TimeoutStopSec=30
StartLimitBurst=3
StartLimitIntervalSec=60
Environment="ENVIRONMENT=production"
Environment="CLERK_SECRET_KEY=sk_test_NVaHy5Dm9iS9dCfn6Fxv2MZxR8aTu0LhNeaxKAyC70"

# Isolation
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/flashback-restore-monorepo
ReadOnlyPaths=/opt/flashback-restore-monorepo/backend/.venv

[Install]
WantedBy=multi-user.target
